# launchd/com.smolpaws.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.smolpaws</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{PROJECT_ROOT}}/scripts/run-local-smolpaws.sh</string>
        <string>npm</string>
        <string>run</string>
        <string>dev:start</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{PATH}}</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>SMOLPAWS_HOME_DIR</key>
        <string>{{SMOLPAWS_HOME_DIR}}</string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{LOG_DIR}}/smolpaws.launchagent.log</string>
    <key>StandardErrorPath</key>
    <string>{{LOG_DIR}}/smolpaws.launchagent.error.log</string>
</dict>
</plist>
